# images/locust/local.conf
[supervisord]
nodaemon=false

[program:cron]
command=/etc/init.d/cron start
user=root
autorestart=false
startretries=0

[program:sysstat]
command=/etc/init.d/sysstat start
user=root
autorestart=false
startretries=0

[program:influxdb]
command=/usr/bin/influxd -pidfile /var/run/influxdb/influxd.pid -config /etc/influxdb/influxdb.conf

[program:telegraf]
command=/usr/bin/telegraf -config /etc/telegraf/telegraf.conf -config-directory /etc/telegraf/telegraf.d

[program:locust]
command=/bin/bash -c "locust -f LOCUST_SCRIPT --host=WEB_APP_ADDRESS \"$@\""
user=root
autorestart=false
startretries=0